# wayvnc/.config/systemd/user/wayvnc.service
[Unit]
Description=WayVNC Remote Desktop Server
Documentation=man:wayvnc(1)
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/wayvnc --log-level=trace -C %h/.config/wayvnc/config
Restart=on-failure
RestartSec=5

# Security hardening (minimal set that allows PAM authentication)
# Note: More restrictive settings prevent PAM's unix_chkpwd from accessing /etc/shadow
# NoNewPrivileges=false is required to allow unix_chkpwd (setuid) to read /etc/shadow
NoNewPrivileges=false
PrivateTmp=true

[Install]
WantedBy=graphical-session.target
